# launchd/io.github.bunyodrafikov.agent-maxprime.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.github.bunyodrafikov.agent-maxprime</string>

    <!-- ┌─ EDIT HERE: which AI CLIs to run ─────────────────────────────┐
         │ The first two <string> entries are fixed (bash + runner).     │
         │ Everything after is passed as an agent name to run-checkin.sh │
         │ in order. Add, remove, or reorder the trailing entries.       │
         └───────────────────────────────────────────────────────────────┘ -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO__/bin/run-checkin.sh</string>
        <string>codex</string>
        <string>claude</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- ┌─ EDIT HERE: when to fire ─────────────────────────────────────┐
         │ One <dict> per fire time, 24-hour clock. Add or remove dicts  │
         │ to change frequency. Optional keys per dict: Weekday (0=Sun), │
         │ Day (1-31), Month (1-12).                                     │
         └───────────────────────────────────────────────────────────────┘ -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key><integer>6</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>11</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>16</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Caches/agent-maxprime/launchd.out</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Caches/agent-maxprime/launchd.err</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
